# .config/autostart/xrdb.desktop
[Desktop Entry]
Type=Application
Name=xrdb
Exec=/usr/bin/xrdb -merge -I$HOME .Xresources
StartupNotify=false
Terminal=false
