# etc/udev/rules.d/89-wifi-poweroff.rules
# disable wifi power save
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wl*", RUN+="/usr/bin/iwconfig wlp59s0 power off"
